# deploy/formfruta.service.example
[Unit]
Description=Form Fruta Comercial Streamlit
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=formfruta
Group=formfruta
WorkingDirectory=/opt/form-fruta-comercial
Environment=PYTHONUNBUFFERED=1
Environment=FORM_FRUTA_SOURCE_SYSTEM=form_fruta_comercial_prod_ol9
ExecStart=/opt/form-fruta-comercial/.venv/bin/python -m streamlit run streamlit_app.py --server.address=127.0.0.1 --server.port=8501 --browser.serverAddress=formfruta.intra.empresa.local
Restart=always
RestartSec=5
TimeoutStopSec=20
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/opt/form-fruta-comercial/data
UMask=0077

[Install]
WantedBy=multi-user.target
