# deploy/formfruta-email.timer.example
[Unit]
Description=Form Fruta Comercial periodic email dispatch

[Timer]
OnCalendar=*:0/5
Persistent=true
Unit=formfruta-email.service

[Install]
WantedBy=timers.target
