# deploy/formfruta.service.example
[Unit]
Description=Form Fruta Comercial Docker Stack
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=formfruta
Group=formfruta
WorkingDirectory=/home/soporte/FromFrutaComercial
Environment=COMPOSE_PROJECT_NAME=formfruta
ExecStart=/usr/bin/docker compose -f /home/soporte/FromFrutaComercial/compose.prod.yml up -d --build app
ExecStop=/usr/bin/docker compose -f /home/soporte/FromFrutaComercial/compose.prod.yml stop app
TimeoutStartSec=0
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
